# systemd/backend.service
[Unit]
Description=Node Backend
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/backend
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/backend/.env

[Install]
WantedBy=multi-user.target
